# supervisor.conf
[supervisord]
nodaemon=false
user=root
logfile=/var/log/musicbot/supervisord.log
pidfile=/var/run/supervisord.pid

[program:musicbot]
command=python -m bot
directory=/app
user=root
autostart=true
autorestart=true
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/musicbot/bot.out.log
stderr_logfile=/var/log/musicbot/bot.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=PYTHONUNBUFFERED="1"

# Restart policy - max 3 restarts per 60 seconds
# After that, wait 60 seconds before trying again
exitcodes=0,2
